# etc/deliverer.cron.example
*/10 * * * * training /bin/bash /opt/training/misc/deliverer.sh
