# etc/training.delivery.timer.example
[Unit]
Description=Runs Training Delivery Agent regularly

[Timer]
OnUnitInactiveSec=1800
Unit=training.delivery.service

[Install]
WantedBy=timers.target
