# etc/training.delivery.timer.example
[Unit]
Description=Training Delivery Trigger

[Timer]
OnStartupSec=60
OnUnitInactiveSec=600
Unit=training.delivery.service

[Install]
WantedBy=timers.target
