# etc/training.delivery.service.example
[Unit]
Description=Training Delivery Agent

[Service]
Type=simple
User=training
ExecStart=/opt/training/deliverer.py
